# supervisor/php-fpm.conf
[program:php-fpm]
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log
command=php5-fpm -F